# ansible/files/etc/profile.d/Z999-linuxlogo.sh
# ansible managed file - do not edit manually as changes will be overwritten

neofetch
